# update.timer
[Unit]
Description=update-irc

[Timer]
OnUnitActiveSec=15 min
OnBootSec=20s

[Install]
WantedBy=timers.target
